# deploy/kamir-gpio-play.service
[Unit]
Description=Kamir GPIO Play Session
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/dev/kamir
ExecStart=/home/pi/.local/bin/kamir gpio-play
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
