# deploy/kamir-gpio-play.service
[Unit]
Description=Kamir GPIO Play Session
After=network.target

[Service]
Type=simple
User=pi
# %h expands to the home directory of User= above (/home/pi).
# Adjust the repo path (dev/kamir) if you cloned elsewhere.
WorkingDirectory=%h/dev/kamir
ExecStart=%h/.local/bin/kamir gpio-play
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
